# ansible/templates/gunicorn.service.j2
[Unit]
Description=URL Shortener Gunicorn daemon
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
EnvironmentFile={{ app_dir }}/.env
ExecStart={{ venv_dir }}/bin/gunicorn \
    --workers 2 \
    --bind 0.0.0.0:8000 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    url_shortener_app.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
